# deploy/docker-proxy.service
[Unit]
Description=Docker Registry Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/docker-proxy
ExecStart=/opt/docker-proxy/docker-proxy
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
